# root/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

echo 35 > /sys/devices/system/cpu/cpufreq/ondemand/up_threshold
echo 10 > /sys/devices/system/cpu/cpufreq/ondemand/sampling_down_factor

exit 0
